# roles/ipa_backup/templates/ipa-backup-iriven.service.j2
[Unit]
Description=IRIVEN FreeIPA backup job
Documentation=man:ipa-backup(1)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/sbin/ipa-backup-iriven.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
PrivateTmp=true
NoNewPrivileges=true
