# bkp/roles/ipa_monitoring_baseline/templates/idm-monitoring-check.service.j2
[Unit]
Description=IRIVEN IDM monitoring baseline checks
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart={{ idm_monitoring_scripts_dir }}/idm-monitoring-check.sh
Nice=10
PrivateTmp=true
NoNewPrivileges=true
