# roles/ipa_monitoring_baseline/templates/idm-monitoring-check.timer.j2
[Unit]
Description=Run IRIVEN IDM monitoring baseline checks

[Timer]
OnCalendar={{ idm_monitoring_timer_on_calendar }}
Persistent=true
Unit=idm-monitoring-check.service

[Install]
WantedBy=timers.target
