# server/named.service.d/override.conf
[Service]
Type=simple
PIDFile=
ExecStart=
ExecStart=/usr/sbin/named -u named -c ${NAMEDCONF} $OPTIONS -g
TimeoutStartSec=180
